# log2ram-daily.service
[Unit]
Description=Daily Log2Ram writing activities
After=log2ram.service

[Service]
ExecStart=/bin/systemctl reload log2ram.service
