# log2ram.cron
#!/usr/bin/env sh

systemctl reload log2ram
