# debian/postinst
#!/usr/bin/env bash
set -euo pipefail

systemctl enable log2ram
echo "#####         Reboot to activate log2ram         #####"
echo "##### edit /etc/log2ram.conf to configure options ####"